# smartsort.service
[Unit]
Description=SmartSort File Organizer Service
After=network.target

[Service]
Type=simple
# User-specific path - Change these paths when deploying to another machine if your installation folder is different
# Example: Change %h/project/smartsort to your actual project directory (e.g., /home/username/smartsort)
ExecStart=%h/project/smartsort/smartsort/bin/python %h/project/smartsort/main.py --daemon
Restart=on-failure

[Install]
WantedBy=default.target
